# ppe_safety_gating/ppe-safety-gating.service
[Unit]
Description=PPE Safety Gating System - Hard Safety Enforcement
After=network.target
StartLimitBurst=5
StartLimitIntervalSec=30s

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ppe_safety_gating

# Relay mode example:
ExecStart=/usr/bin/python3 main.py \
  --rtsp rtsp://192.168.0.5:554/stream \
  --model yolov8n.pt \
  --mode relay \
  --relay-pin 17 \
  --green-pin 27 \
  --red-pin 22 \
  --buzzer-pin 23 \
  --db /var/log/ppe_safety_events.db

# For PLC mode, use:
# ExecStart=/usr/bin/python3 main.py \
#   --rtsp rtsp://192.168.0.5:554/stream \
#   --model yolov8n.pt \
#   --mode plc \
#   --plc-host 192.168.0.10 \
#   --plc-port 502 \
#   --plc-coil 1 \
#   --green-pin 27 \
#   --red-pin 22 \
#   --buzzer-pin 23 \
#   --db /var/log/ppe_safety_events.db

# Restart policy: fail-safe behavior
Restart=on-failure
RestartSec=10s
MaxStartCount=5

# Enable core dumps for debugging crashes
LimitCORE=infinity

# Logging to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ppe-safety-gating

# Security hardening (optional)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
